# rootfs/humhub.cron
* * * * * /usr/local/bin/php /var/www/html/protected/yii queue/run >/dev/null 2>/dev/stderr
* * * * * /usr/local/bin/php /var/www/html/protected/yii cron/run >/dev/null 2>/dev/stderr
