# rootfs/humhub.cron
* * * * * /usr/local/bin/php /var/www/html/protected/yii queue/run
* * * * * /usr/local/bin/php /var/www/html/protected/yii cron/run
